# home/systemd/pacman-synchronize.service
[Unit]
Description=Pacman - synchronize package database

[Service]
Type=simple
ExecStart=/usr/bin/pacman -Sy

